# resources/etc/supervisor/conf.d/curl.conf
[program:curl]
command=curl
autostart=true
autorestart=true
stderr_logfile=/var/log/curl.err.log
stdout_logfile=/var/log/curl.out.log
