# resources/etc/supervisor/conf.d/maven.conf
[program:maven]
command=mvn
autostart=true
autorestart=true
stderr_logfile=/var/log/maven.err.log
stdout_logfile=/var/log/maven.out.log
