# resources/etc/supervisor/conf.d/pip.conf
[program:pip]
command=pip
autostart=true
autorestart=true
stderr_logfile=/var/log/pip.err.log
stdout_logfile=/var/log/pip.out.log
