# resources/etc/supervisor/conf.d/python3.conf
[program:python3]
command=python3
autostart=true
autorestart=true
stderr_logfile=/var/log/python3.err.log
stdout_logfile=/var/log/python3.out.log
